# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/run/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:server]
command=/app/server
directory=/app
autostart=true
autorestart=true
startretries=3
priority=20
stdout_logfile=/app/logs/server.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/server.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
